# docker/rootfs/etc/s6-overlay/s6-rc.d/cron/run
#!/command/with-contenv bash
if [ "${WEBTLO_CRON:=false}" = "true" ]; then
  s6-echo "Enabled cron jobs"
  exec crond -f
else
  s6-echo "Running with disabled cron"
  exec sleep infinity
fi
