# docker/rootfs/etc/s6-overlay/s6-rc.d/php-fpm/run
#!/bin/sh

/usr/sbin/php-fpm82 --nodaemonize
