# docker/rootfs/etc/s6-overlay/s6-rc.d/php-fpm/run
#!/command/execlineb -P
php-fpm81 --nodaemonize
