# docker/rootfs/etc/s6-overlay/s6-rc.d/_setup-storage/up
#!/command/execlineb
foreground { bash -e /root/_setup-storage.sh }
